# update.coreelec.org/crontab
0 0 * * * /usr/bin/php /your/path/to/clean-statsdb.php > /dev/null 2>&1
